# roles/baseline/templates/etc/profile.d/baseline.sh
alias splkuser='sudo -iu splunk'
alias splkcom='sudo -su splunk'
alias splksys='sudo su -ls /bin/bash splunk'
